# systemd/doorbell.service
[Unit]
Description=玄関モニタリングシステム (Doorbell Monitor)
After=network.target
Wants=network-online.target

[Service]
# 実行ユーザー (通常は "pi" または自分のユーザー名に変更してください)
User=toya
Group=toya

# プロジェクトディレクトリ (実際のパスに変更してください)
WorkingDirectory=/home/toya/LineBot

# 仮想環境の Python を使用
ExecStart=/home/toya/LineBot/venv/bin/python /home/toya/LineBot/main.py

# 異常終了時に自動再起動
Restart=always
RestartSec=10

# ログをジャーナルに出力
StandardOutput=append:/home/toya/LineBot/doorbell.log
StandardError=append:/home/toya/LineBot/doorbell.log

# GPIO アクセスのための環境変数
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
